# root/etc/s6-overlay/s6-rc.d/svc-mod-deluge-dynamic-incoming-port/run
#!/usr/bin/with-contenv bash

cmd=(exec s6-setuidgid abc
  /usr/bin/deluge-dynamic-incoming-port
)

if [[ -n "${DELUGE_DYNAMIC_INCOMING_PORT_DEBUG:-}" ]]; then
  cmd+=( "--verbose" )
fi

"${cmd[@]}"
