# root/etc/s6-overlay/s6-rc.d/svc-mod-deglue-dynamic-incoming-port/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc run /usr/bin/deluge-dynamic-incoming-port
